# deploy/laptop/tilbudsfinder-scraper-agent-update.timer
[Unit]
Description=Check for Tilbudsfinder scraper agent image updates every five minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
